# henhouse-reinit.service
[Service]
Type=oneshot
ExecStart=/bin/sh -c "echo EXTRA=--reinit > /var/run/henhouse/extra"
ExecStart=/bin/systemctl restart henhouse
