# henhouse.service
[Unit]
Description=Scoreboard for jeopardy-style CTF
Requires=postgresql.service
After=postgresql.service

[Service]
User=henhouse
Group=henhouse
WorkingDirectory=/var/lib/henhouse/
ExecStartPre=/bin/mkdir -p /var/run/henhouse/
ExecStart=/bin/bash -c "/usr/bin/henhouse /etc/henhouse.toml ${EXTRA}"
ExecStartPost=/bin/rm -f /run/henhouse/extra
EnvironmentFile=/run/henhouse/extra

[Install]
WantedBy=multi-user.target
